# dot_config/systemd/user/arbtt.service
[Unit]
Description=arbtt time tracking

[Service]
ExecStart=/usr/bin/arbtt-capture
Restart=always

[Install]
WantedBy=default.target